# screenshot-api.service
[Unit]
Description=Screenshot & PDF API
After=network.target

[Service]
Type=simple
User=agent
WorkingDirectory=/home/agent/projects/screenshot-api
ExecStart=/usr/bin/node server.js
Restart=always
RestartSec=5
Environment=PORT=3001
Environment=NODE_ENV=production

[Install]
WantedBy=multi-user.target
